# bash/profile.d/aliases
alias vi="vim"
#alias ls="ls --color=auto"
#alias ll="ls --color -lh"
#alias lll="ls --color -lha"
#alias la="ls --color -la"
alias grep="grep --color"
alias pgrep="pgrep -f"
alias pkill="pkill -f"
alias egrep="egrep --color"
alias fgrep="fgrep --color"
alias cls="clear"
alias cd..="cd .."
alias rm="rm -i"
alias mv="mv -i"
#alias cp="cp -i"
alias ln="ln -i"
alias wget='wget -c'
alias df='df -h'
alias du='du -h'
alias free='free -m'
alias sshfs='sshfs -o follow_symlinks'

command -v colordiff >/dev/null 2>&1 && alias diff='colordiff'  # install colordiff package :)
alias path='echo -e ${PATH//:/\\n}'
alias ping='ping -c 5'
alias ports='netstat -tulanp tcp'
alias bc='bc -l' # run with standard math lib

## Mac OS X dependent settings
if [ $(uname -s) = "Darwin" ]; then
    ## aliases for mysql installed from mysql.com 
	if [ -d /usr/local/mysql ]; then
		alias mysql=/usr/local/mysql/bin/mysql
		alias mysqladmin=/usr/local/mysql/bin/mysqladmin
	fi

    ## Vim crontab fix
    alias crontab="VIM_CRONTAB=true crontab"
    
    # tree command
    alias tree="find . -print | sed -e 's;[^/]*/;|____;g;s;____|; |;g'"
fi

## application aliases
# subversion
alias svn_extract_dot_files='find . -name .svn -print0 | xargs -0 rm -rf'
alias svnadd='svn st | grep "^\?" | awk "{print \$2}" | xargs svn add'
alias svndiff='svn diff --diff-cmd ~/bin/svnvimdiff'

#additional alias
alias chmox='chmod +x'

##git
# pretty log
alias gitlg="git log --graph --pretty=format:'%Cred%h%Creset - %s %Cgreen(%cr)%Creset%C(yellow)%d%Creset %C(dim white)%an%Creset' --abbrev-commit --date=relative --date-order --all" 

#detailed commit logs
alias gitcilg='git log --pretty=format:"[%h] %ae, %ar: %s" --stat'

#clean screen
alias cleanscreen="screen -ls | tail -n +2 | head -n -2 | awk '{print $1}'| xargs -Ixx screen -S xx -X quit"

# pastebin-cl private paste alias
alias pastebincl='pastebincl -p'
