# bash/profile.d/mybashrc
#PS1='\u@\h:\W \$ '
PS1='\[\033[1m\]\u@\h \[\033[0;31m\]\W \$ \[\033[00m\]'

# make bash history unlimited
export HISTFILESIZE=
export HISTSIZE=

# set editor as vim
export EDITOR="vim"
export PAGER="less"

#vi-compatiple command line
#set -o vi

# Language preferences
export LANG="en_US.UTF-8"
export LC_ALL="en_US.UTF-8"

# Terminal colors
export CLICOLOR=1
## White background
#export LSCOLORS=ExFxCxDxBxegedabagacad
## Dark background
export LSCOLORS=gxBxhxDxfxhxhxhxhxcxcx


## Add home local lib and bin to the PATH
if [ -e "$HOME/.local/lib" ]; then
    export LD_LIBRARY_PATH="$HOME/.local/lib:$LD_LIBRARY_PATH"
fi

if [ -e "$HOME/.local/lib64" ]; then
    export LD_LIBRARY_PATH="$HOME/.local/lib64:$LD_LIBRARY_PATH"
fi

if [ -d "$HOME/.local/bin" ]; then
    export PATH="$HOME/.local/bin:$PATH"
fi

## Mac OS X dependent settings
if [ $(uname -s) = "Darwin" ]; then

    # Java home
    export JAVA_HOME=$(/usr/libexec/java_home)

    # Autojump for bash installed with macports
    if [ -f /usr/local/etc/profile.d/autojump.sh ]; then
    source /usr/local/etc/profile.d/autojump.sh
        export AUTOJUMP_IGNORE_CASE=1
        export AUTOJUMP_KEEP_SYMLINKS=1
    fi

    # Mac only
    export BASH_SILENCE_DEPRECATION_WARNING=1
fi


## Other
if [ $(uname -s) = "Linux" ]; then

    # autojump for bash (ubuntu package)
    if [ -f /usr/share/autojump/autojump.bash ]; then
        . /usr/share/autojump/autojump.bash
        export AUTOJUMP_IGNORE_CASE=1
        export AUTOJUMP_KEEP_SYMLINKS=1
    fi

    ## On Arch Linux it is installed under /usr/etc/
    if [ -f /usr/etc/profile.d/autojump.bash ]; then
        . /usr/etc/profile.d/autojump.bash
        export AUTOJUMP_IGNORE_CASE=1
        export AUTOJUMP_KEEP_SYMLINKS=1
    fi

    #if reboot required tell me (Ubuntu)
    if [ -f /var/run/reboot-required ]; then
        echo 'Reboot required'
    fi
fi

# Explicitly set HOSTNAME variable
export HOSTNAME=`hostname`

# Setup direnv if exists
if [ -x "$(command -v direnv)" ]; then
    eval "$(direnv hook bash)"
fi
